# launchd/com.xxljob.boost.executor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.xxljob.boost.executor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/song/.lark-channel-workspaces/codex/default/xxl-job-boost/scripts/run-executor.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/song/.lark-channel-workspaces/codex/default/xxl-job-boost</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/xxl-job-boost-run/executor.launchd.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/xxl-job-boost-run/executor.launchd.err</string>
</dict>
</plist>
