# launchd/com.codex.continue-20260609-0930.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.codex.continue-20260609-0930</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/song/.lark-channel-workspaces/codex/default/xxl-job-boost/scripts/remind-continue-20260609-0930.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/song/.lark-channel-workspaces/codex/default/xxl-job-boost</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Year</key>
    <integer>2026</integer>
    <key>Month</key>
    <integer>6</integer>
    <key>Day</key>
    <integer>9</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/xxl-job-boost-run/continue-20260609-0930.launchd.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/xxl-job-boost-run/continue-20260609-0930.launchd.err</string>
</dict>
</plist>
